# conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autostart=true
autorestart=true
# stdout_logfile=/var/log/fluxbox.log
# stderr_logfile=/var/log/fluxbox.err.log