# conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autostart=true
autorestart=true
