# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display :${VNC_DISPLAY}.${VNC_SCREEN} -forever -usepw -shared -rfbport ${VNC_PORT}
autostart=true
autorestart=true
