# conf.d/xterm.conf
[program:xterm]
command=/usr/bin/xterm
autostart=%(ENV_AUTO_START_XTERM)s
autorestart=true
stdout_logfile=/var/log/xterm.log
stderr_logfile=/var/log/xterm.err.log
